# baseline/airootfs/etc/xdg/autostart/tealinux-installer.desktop
[Desktop Entry]
Name=Tealinux Installer
Categories=Development
Exec=tealinux-installer
Icon=tealinux-installer
Terminal=false
Type=Application
